# files/supervisord.conf
[program:rethinkdb]
command=rethinkdb serve --config-file /etc/rethinkdb/instances.d/default.conf
